# deploy/operations/systemd/gray-mysql-restore-verify.timer
[Unit]
Description=Verify the latest Gray Commerce MySQL backup weekly

[Timer]
OnCalendar=Sun *-*-* 04:30:00
Persistent=true
RandomizedDelaySec=30m
Unit=gray-mysql-restore-verify.service

[Install]
WantedBy=timers.target
